# net-wireless/eiwd/files/ead.initd
#!/sbin/openrc-run
# Copyright 1999-2023 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

pidfile="/run/ead.pid"
command="/usr/libexec/ead"
command_background="yes"

depend() {
	keyword -shutdown
}
